# Linux/reconnect.service
[Unit]
Description=Auto reconnect to network

[Service]
Type=simple
ExecStart=/usr/bin/reconnect

[Install]
WantedBy=multi-user.target
